# crontab
# crontab
# Define las tareas programadas. La salida se redirige a los logs de Docker/Railway.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# A las 12:00 UTC todos los días, ejecuta el pipeline de scraping y análisis.
# El `&&` asegura que el analizador solo corra si el scraper tiene éxito.
0 12 * * * /usr/local/bin/python3 /app/scraper.py >> /proc/1/fd/1 2>/proc/1/fd/2 && /usr/local/bin/python3 /app/analyzer.py >> /proc/1/fd/1 2>/proc/1/fd/2

# A las 14:00 UTC todos los días, ejecuta el script de monitoreo.
0 14 * * * /usr/local/bin/python3 /app/monitor.py >> /proc/1/fd/1 2>/proc/1/fd/2

# Es una buena práctica dejar una línea en blanco al final del archivo crontab.